# dashboarddeltas/dashboarddeltas.conf
# checks for changes to COVID-19 website
[program:dashboarddeltas]
user=ubuntu
directory=/home/ubuntu/Development/covid19-alerts/dashboarddeltas
command=python3 deltabot.py
autostart=true
autorestart=true
startsecs=5
